# root/etc/services.d/celery/run
#!/usr/bin/with-contenv sh

cd /app && /usr/local/bin/celery -A internal worker -Q celery -B -l info --hostname=docker_internal